# demo/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:rackup]
command=/bin/bash -c "/opt/demo/bin/start-server.sh"